# roles/autoupdate/templates/podman-prune.service.j2
# {{ ansible_managed }}
[Unit]
Description=Scheduled Podman images update

[Service]
Type=simple
ExecStart=/usr/bin/podman auto-update